# templates/default/50-ethtool.rules.erb
# Generated by Chef for <%= node['fqdn'] %>
# Local modifications will be overwritten.
<% node['iscsi']['udev']['interfaces'].each do |iscsi_if| %>
  <% node['iscsi']['udev']['ethtool_opts'].each do |ethtool_opt| %>
ACTION=="add", SUBSYSTEM=="net", KERNEL=="<%= iscsi_if %>", RUN+="/sbin/ethtool <%= ethtool_opt %>"
  <% end %>
<% end %>
