# templates/default/50-ethtool.rules.erb
# Generated by Chef for <%= node["fqdn"] %>
# Local modifications will be overwritten.
<% node["iscsi"]["interfaces"].each do |iscsi_if| %>
<% node["iscsi"]["ethtool_opts"].each do |ethtool_opt| %>
ACTION=="add", SUBSYSTEM=="net", NAME=="<%= iscsi_if %>", RUN+="/sbin/ethtool <%= ethtool_opt %>"
<% end %>
<% end %>
